# backup-server.timer
[Unit]
Description=Backup server data

[Timer]
OnUnitActiveSec=60m

[Install]
WantedBy=timers.target
