# update-server.timer
[Unit]
Description=Backup server data

[Timer]
OnCalendar=5:0:00

[Install]
WantedBy=timers.target
