# update-server.timer
[Unit]
Description=Backup server data

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
